# plugins/dev/scripts/orch-monitor/dist/ai.coalesce.catalyst-health-responder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the periodic cloud-sync health responder (CTL-1509).

  Do not edit this file directly — tokens are substituted by install-health-responder.sh.

  Instructions:
    Run: bash plugins/dev/scripts/install-health-responder.sh
    (or it is called automatically by 'catalyst-stack install-services')

  Manual steps (if needed):
    1. Copy to ~/Library/LaunchAgents/ai.coalesce.catalyst-health-responder.plist
    2. Replace REPLACE_WITH_ABSOLUTE with the absolute path to the scripts dir
       (e.g. /Users/you/catalyst/plugin-source/plugins/dev/scripts)
    3. Replace REPLACE_HOME with your home directory (e.g. /Users/you)
    4. Replace REPLACE_START_INTERVAL with the interval in seconds (e.g. 180)
    5. Load:
         launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/ai.coalesce.catalyst-health-responder.plist
    6. Verify:
         launchctl list | grep health-responder
         tail -f ~/catalyst/health-responder.log

  To unload:
    bash plugins/dev/scripts/install-health-responder.sh --uninstall

  For a one-shot test run:
    bash REPLACE_WITH_ABSOLUTE/health-responder.sh --dry-run

  See: docs/health-responder.md
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.coalesce.catalyst-health-responder</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPLACE_WITH_ABSOLUTE/health-responder.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>REPLACE_START_INTERVAL</integer>

  <key>RunAtLoad</key>
  <true/>

  <!-- launchd's default PATH omits /opt/homebrew/bin — without it the sweep
       cannot resolve jq and every self-heal breadcrumb is silently ignored
       (no settling hold, no no-respawn detection). Substituted by the
       installer; mirrors catalyst-stack's _stack_agent_path (CTL-1289). -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>REPLACE_PATH</string>
    <!-- launchd does NOT inherit the installing shell's environment: a node
         using a nondefault CATALYST_DIR must have it PERSISTED here or the
         installed sweep would watch $HOME/catalyst — the wrong lock,
         breadcrumb, and marker dir. Substituted at install time. -->
    <key>CATALYST_DIR</key>
    <string>REPLACE_CATALYST_DIR</string>
  </dict>

  <key>StandardOutPath</key>
  <string>REPLACE_HOME/catalyst/health-responder.log</string>

  <key>StandardErrorPath</key>
  <string>REPLACE_HOME/catalyst/health-responder.log</string>
</dict>
</plist>
