# plugins/dev/scripts/channel-watcher/ai.coalesce.catalyst-channel-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.coalesce.catalyst-channel-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPLACE_WITH_LAUNCH_SH_ABS_PATH</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>REPLACE_WITH_PATH</string>
    <key>HOME</key>
    <string>REPLACE_WITH_HOME</string>
    <key>CATALYST_DIR</key>
    <string>REPLACE_WITH_HOME/catalyst</string>
  </dict>

  <key>StandardOutPath</key>
  <string>REPLACE_WITH_HOME/catalyst/channel-watcher.log</string>

  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_HOME/catalyst/channel-watcher.log</string>
</dict>
</plist>
